# packages/chalmers/chalmers-0.7.2.tar.gz/chalmers-0.7.2/chalmers/service/data/upstart.conf
# Chalmers

description "Chalmers process manager"
author "Continuum Analytics <support@continuum.io>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]


start on startup

# Start the process
exec {launch} -c "{python_exe} {chalmers} start --all"

post-start script
   logger -t chalmers "upstart: post-start: job=$UPSTART_JOB"
end script
