# packages/peekaboo/peekaboo-0.6-5.tar.gz/peekaboo-0.6-5/contrib/peekaboo.service
[Unit]
Description=Peekaboo
After=network.target

[Service]
Restart=on-failure
RestartSec=10
Type=simple
User=peekaboo
Group=peekaboo
ExecStart=/usr/bin/peekaboo

[Install]
WantedBy=multi-user.target
