# packages/rattail/rattail-0.4.4.tar.gz/rattail-0.4.4/data/etc/init.d/rattail-dbsync_tmpl
#! /bin/sh
### BEGIN INIT INFO
# Provides:          %(name)s
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: %(desc)s
# Description:       Manages the %(desc)s daemon.
### END INIT INFO

DESC="%(desc)s"
NAME=%(name)s

PYTHON=%(envdir)s/bin/python
RATTAIL=%(envdir)s/bin/rattail
RATTAIL_ARGS="--config=%(envdir)s/app/dbsync.ini"

. /var/lib/rattail/init.d/dbsync
