# packages/Elpotrero/Elpotrero-1.6.2.tar.gz/Elpotrero-1.6.2/elpotrero/_files/tree/scripts/conf/supervisor/_files/PROJECT.conf
[program:{{project}}.{{branch}}]
user=root
command=bash {{rootpath}}/{{project}}.{{domain}}/scripts/gunicorn.sh
stdout_logfile={{rootpath}}/{{project}}.{{domain}}/logs/supervisor/out.log
stderr_logfile={{rootpath}}/{{project}}.{{domain}}/logs/supervisor/err.log
