# packages/easy-aws/easy-aws-1.0.6.tar.gz/easy-aws-1.0.6/easyaws/templates/supervisord.conf
[program:${PROJECT_NAME}]
command=/home/ubuntu/virtualenv/bin/gunicorn ${PROJECT_NAME}.wsgi:application
directory=${PROJECT_PATH}
user=www-data
autostart=true
autorestart=true
redirect_stderr=true